# archlive/airootfs/etc/profile.d/expand_cowsize.sh
#!/bin/bash
mount -o remount,size=6G /run/archiso/cowspace